# debian/gemian-setup.service
[Unit]
Description=Gemian Setup
Before=local-fs.target var-lib-lxc-android-rootfs.mount

[Service]
Type=oneshot
ExecStart=+/usr/sbin/gemian-setup
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
